# nvidia_jetson/ROOTFS/etc/systemd/system/dnsmasq.service.d/override.conf
[Unit]
# start after the lan0 bound interface is up
Requires=network-online.target
After=network-online.target

[Service]
Restart=always
RestartSec=30
